# 车辆贷款违约预测（讯飞A.I算法赛）/deploy/vm/supervisor_flume.conf
# ============================================================
#  Supervisor Configuration — 贷款客户信息修复与智能决策系统
# ============================================================
# Supervisor 用于管理长时间运行的进程（Flume、Kafka、Spark Streaming 等），
# 提供进程保活、自动重启和日志管理。
#
# 安装 Supervisor：
#   sudo apt-get install -y supervisor
#
# 部署配置：
#   sudo cp deploy/vm/supervisor_flume.conf /etc/supervisor/conf.d/loan-flume.conf
#   sudo systemctl restart supervisor
#
# 查看状态：
#   supervisorctl status
#   supervisorctl tail loan-flume
# ============================================================

[program:loan-flume]
command=/opt/flume/current/bin/flume-ng agent -n loan-agent -c /opt/flume/conf -f /opt/flume/conf/flume-agent.conf
directory=/opt/flume/current
autostart=true
autorestart=true
startretries=3
exitcodes=0
user=root
redirect_stderr=true
stdout_logfile=/data/flume/logs/supervisor_loan_flume.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/data/flume/logs/supervisor_loan_flume_err.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
environment=PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[program:loan-kafka]
command=/opt/kafka/kafka/bin/kafka-server-start.sh /opt/kafka/kafka/config/kraft/server.properties
directory=/opt/kafka/kafka
autostart=true
autorestart=true
startretries=3
exitcodes=0
user=root
redirect_stderr=true
stdout_logfile=/data/kafka/logs/supervisor_loan_kafka.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stderr_logfile=/data/kafka/logs/supervisor_loan_kafka_err.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
environment=KAFKA_OPTS="-Xmx2g -Xms2g",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[program:loan-spark-streaming]
command=/opt/bigdata/spark/bin/spark-submit --master yarn --deploy-mode cluster --executor-memory 2g --executor-cores 2 --conf spark.streaming.backpressure.enabled=true /data/loan/jobs/streaming/run_streaming.sh --mode kafka
directory=/data/loan
autostart=false
autorestart=true
startretries=3
exitcodes=0,137
user=root
redirect_stderr=true
stdout_logfile=/data/loan/logs/supervisor_spark_streaming.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stderr_logfile=/data/loan/logs/supervisor_spark_streaming_err.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
environment=HADOOP_HOME="/opt/bigdata/hadoop",SPARK_HOME="/opt/bigdata/spark",JAVA_HOME="/usr/lib/jvm/java-11-openjdk-amd64",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[program:loan-flask-api]
command=/opt/bigdata/loan-venv/bin/python -m flask run --host 0.0.0.0 --port 5000
directory=/data/loan
autostart=true
autorestart=true
startretries=5
exitcodes=0
user=root
redirect_stderr=true
stdout_logfile=/data/loan/logs/supervisor_flask.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/data/loan/logs/supervisor_flask_err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3
environment=FLASK_APP="service/flask/app.py",PYTHONPATH="/data/loan",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
